# crontab
* * * * * echo "hello 1 Minute"
*/5 * * * * echo "hello 5 Minute"
# remember to end this file with an empty new line
