# src/crontab
* * * * * cd /var/www/html/src/; php /var/www/html/src/time.php -q -f
*/5 * * * * echo "hello 5 Minute"
# remember to end this file with an empty new line
